# packages/pywizard-extra-modules/pywizard-extra-modules-0.1.15.tar.gz/pywizard-extra-modules-0.1.15/pywizard_extra_modules/templates/nodejs-na-live.conf
#!upstart
description "customer support server"
author      "Sergei Vizel"

env RUNAS=www-data
env APP=app.js
env NODEJS="NODE_ENV=production node"
env HOMEDIR="/home/alex/dev/app/src/MFA/CustomerSupportBundle/Resources/nodejs"
 
start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown
 
script
    chdir $HOMEDIR
    exec sudo -u $RUNAS $NODEJS $APP >> /var/log/${UPSTART_JOB}.log 2>&1
end script
 
pre-start script
    # Date format same as (new Date()).toISOString() for consistency
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Starting" >> /var/log/${UPSTART_JOB}.log
end script
 
pre-stop script
    echo "[`date -u +%Y-%m-%dT%T.%3NZ`] (sys) Stopping" >> /var/log/${UPSTART_JOB}.log
end script






