# packages/xsanio-server/xsanio-server-0.5.5.tar.gz/xsanio-server-0.5.5/osx_launchd/ru.shortcut.xsanio_server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ru.shortcut.xsanio_server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/gunicorn_django</string>
        <string>--bind=:80</string>
        <string>--workers</string>
        <string>4</string>
        <string>/opt/xsanio_server/xsanio_server/settings.py</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>UserName</key>
    <string>root</string>
</dict>
</plist>
