# packages/baruwa/baruwa-2.0.0.tar.bz2/baruwa-2.0.0/debian/baruwa.cron.hourly
#!/bin/bash
#
#
indexer auditlog lists domains accounts organizations archive --rotate &>/dev/null
