# packages/bamboo.scaffold/bamboo.scaffold-0.0.7.tar.gz/bamboo.scaffold-0.0.7/bamboo/scaffold/tmpls/upstart-rhaptos2.conf
# Rhaptos2                                                                                 
# ToDO: privileges - run as www-data or other user 
# http://superuser.com/questions/213416/running-upstart-jobs-as-unprivileged-users                                                                                   


description     "Rhaptos2"

start on runlevel [2345]
stop on runlevel [!2345]

#enable starting multiple instances on one server, diff ports
instance $PORT


pre-start script
    echo '**** start on $PORT'    |	  logger -ist 'rhaptos2'					   
end script

post-stop script
    echo '*** stopping $PORT' | logger -ist 'rhaptos2'					   
end script


## From the confd passed to this template,
## set and export the config values this app needs to run
## (I cannot easily see a nicer way in upstart)

env rhaptos2repo_repodir=%(rhaptos2repo_repodir)s
env rhaptos2repo_statsd_host=%(rhaptos2repo_statsd_host)s
env rhaptos2repo_statsd_port=%(rhaptos2repo_statsd_port)s
env rhaptos2repo_www_server_name=%(rhaptos2repo_www_server_name)s
env rhaptos2repo_cdn_server_name=%(rhaptos2repo_cdn_server_name)s
env rhaptos2repo_openid_secretkey=%(rhaptos2repo_openid_secretkey)s
env rhaptos2repo_use_logging=%(rhaptos2repo_use_logging)s
env rhaptos2repo_loglevel=%(rhaptos2repo_loglevel)s


export rhaptos2repo_repodir
export rhaptos2repo_statsd_host
export rhaptos2repo_statsd_port
export rhaptos2repo_www_server_name
export rhaptos2repo_cdn_server_name
export rhaptos2repo_openid_secretkey
export rhaptos2repo_loglevel
export rhaptos2repo_use_logging

exec /usr/bin/python /usr/local/bin/rhaptos2_runrepo.py \
                     --host='0.0.0.0' \
                     --port=$PORT \
                     --debug=True \
         	            2>&1 | logger -ist 'rhaptos2'

respawn
respawn limit 3 20
