# packages/bamboo.scaffold/bamboo.scaffold-0.0.7.tar.gz/bamboo.scaffold-0.0.7/bamboo/scaffold/tmpls/upstart-rhaptos2user.conf
# Rhaptos2                                                                                 
# ToDO: privileges - run as www-data or other user 
# http://superuser.com/questions/213416/running-upstart-jobs-as-unprivileged-users                                                                                   


description     "Rhaptos2User"

start on runlevel [2345]
stop on runlevel [!2345]

#enable starting multiple instances on one server, diff ports
instance $PORT


pre-start script
    echo '**** start on $PORT'    |	  logger -ist 'rhaptos2user'					   
end script

post-stop script
    echo '*** stopping $PORT' | logger -ist 'rhaptos2user'					   
end script


## From the confd passed to this template,
## set and export the config values this app needs to run
## (I cannot easily see a nicer way in upstart)



env bamboo_statsd_host=%(bamboo_statsd_host)s
env bamboo_statsd_port=%(bamboo_statsd_port)s
env rhaptos2user_use_logging=%(rhaptos2user_use_logging)s
env rhaptos2user_loglevel=%(rhaptos2user_loglevel)s

export bamboo_statsd_host
export bamboo_statsd_port
export rhaptos2user_use_logging
export rhaptos2user_loglevel

exec /usr/bin/python /usr/local/bin/run_rhaptos2user.py \
                     --host='0.0.0.0' \
                     --port=$PORT \
                     --debug=True \
         	            2>&1 | logger -ist 'rhaptos2user'

respawn
respawn limit 3 20
